# TODO/etc/udev/rules.d/90-logitechG203.rules
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c084", MODE="660", GROUP="plugdev", RUN+="/usr/bin/g213colors-gui -t"
