# files/etc/systemd/system/zramswap.service.d/override.conf
[Unit]
After=zram-config.service
Wants=zram-config.service
